# cron/crontab
* * * * * /home/fabian/eth-miner/cron/cron.sh >> /home/fabian/eth-miner/cron/script.log 2>&1