# cron/crontab
* * * * * /eth-miner/cron/cron.sh >> /eth-miner/cron/script.log 2>&1